# srcs/mysql/supervisord.conf
[supervisord]
nodaemon = true
user = root

[program:mysql]
command = sh -c "mysqld && kill $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "telegraf; kill $(cat supervisord.pid)"