# srcs/nginx/supervisord.conf
[supervisord]
nodaemon = true
user = root

[program:ssh]
command = sh -c "/usr/sbin/sshd -D && kill $(cat supervisord.pid)"

[program:nginx]
command = sh -c "nginx -g 'daemon off;' && kill $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "telegraf && kill $(cat supervisord.pid)"